# roles/karrot-backend/templates/geoipupdate.service
[Unit]
Description=Update maxmind database service

[Service]
Type=oneshot
ExecStart=/usr/bin/geoipupdate
