# roles/karrot-backend/templates/systemd.worker.target.j2
[Unit]
Description={{ site }} | django workers

{% for n in range(workers) %}
Wants={{ site }}-worker@{{ n+1 }}.service
{% endfor %}

[Install]
WantedBy=multi-user.target